# server/src/main/config/trumpet-server.ini
[program:trumpet-server]
command=/opt/trumpet/server/bin/trumpet.sh --zk.connect zookeeper.server1:2181,zookeeper.server2:2181 --topic trumpet.topic.name --dfs.edits.dir /path/to/hadoop/dfs/jn
autostart=true
autorestart=true
startsecs=10
startretries=999
user=hdfs
stopsignal=TERM
stopasgroup=true
stdout_logfile=/opt/trumpet/server/logs/stdout
stdout_logfile_maxbytes=10000000
stdout_logfile_backups=5
stderr_logfile=/opt/trumpet/server/logs/stderr
stderr_logfile_maxbytes=10000000
stderr_logfile_backups=5
