# server/src/main/config/trumpet-server.ini
[program:trumpet-server]
command=/opt/trumpet/server/bin/trumpet.sh --zk.connect vsccrtc201-brn1.rtc.vrsn.com:2181 --topic vscc.pv.hdfs.inotify.events --dfs.edits.dir /app/dfs/jn
autostart=true
autorestart=true
startsecs=10
startretries=999
user=hdfs
stopsignal=TERM
stopasgroup=true
stdout_logfile=/opt/trumpet/server/logs/stdout
stdout_logfile_maxbytes=10000000
stdout_logfile_backups=5
stderr_logfile=/opt/trumpet/server/logs/stderr
stderr_logfile_maxbytes=10000000
stderr_logfile_backups=5
